# root/etc/services.d/node-consumer/run
#!/usr/bin/with-contenv bash

exec npm run consumer --prefix /var/www
