# root/etc/services.d/node-consumer/run
#!/usr/bin/with-contenv bash

if [[ ! -z "$AMQP_HOST" ]]
then
  exec npm run consumer --prefix /app
else
  exec npm start --prefix /app
fi
